# assets/package/linux/DEBIAN/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable qtet-daemon.service
systemctl start qtet-daemon.service

exit 0
