# assets/package/linux/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop qtet-daemon.service
systemctl disable qtet-daemon.service

exit 0
